# cloud_sync/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/supervisord/supervisord.log
loglevel=debug

[program:cloud_sync]
command=/usr/local/bin/cloud_sync /home/config/cloud_sync.json
startsecs=10